# event-reporter.conf
[program:event-reporter]
command=/code/event-reporter/cmd/event-reporter/event-reporter --interval 2s
autostart=true
autorestart=true