# fake-lepton.conf
[program:fake-lepton]
command=/server/cmd/fake-lepton/fake-lepton
directory=/server/cmd/fake-lepton

autostart=true
autorestart=true