# management-interface.conf
[program:management-interface]
command=/code/management-interface/managementd
autostart=true
autorestart=true