# thermal-uploader.conf
[program:thermal-uploader]
command=/code/thermal-uploader/thermal-uploader
autostart=true
autorestart=true