# thermal-recorder.conf
[program:thermal-recorder]
command=/code/thermal-recorder/cmd/thermal-recorder/thermal-recorder
autostart=true
autorestart=true